# nmos-registration/lib/systemd/system/nmosregistration.service
[Unit]
Description=nmosregistration

[Service]
ExecStart=/usr/bin/nmosregistration
