# nmos-query/lib/systemd/system/nmosquery.service
[Unit]
Description=nmosquery

[Service]
ExecStart=/usr/bin/nmosquery
