# reverse-proxy/etc/init/proxylisting.conf
description	"nmos-proxylisting"
author		"Andrew Bonney <andrew.bonney@bbc.co.uk>"

start on runlevel [2345]
stop on runlevel [016]
respawn

script
	exec sudo /usr/bin/proxylisting
end script
