# mdns-bridge/lib/systemd/system/mdnsbridge.service
[Unit]
Description=mdnsbridge

[Service]
ExecStart=/usr/bin/mdnsbridge
