# reverse-proxy/lib/systemd/system/proxylisting.service
[Unit]
Description=proxylisting

[Service]
ExecStart=/usr/bin/proxylisting
